# deploy/pi/systemd/ids.service
; IDS systemd unit
; Purpose: manage IDS runtime process lifecycle on Raspberry Pi.
; Fit: ensures auto-start and restart behavior for unattended operation.

[Unit]
Description=IDS Player (placeholder)
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/ids/ids.env
ExecStart=/usr/bin/echo "IDS placeholder"
Restart=always

[Install]
WantedBy=multi-user.target
